# systemd/sleep-manager-waker.service
[Unit]
Description=Sleep Manager Waker Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=sleep-manager
Group=sleep-manager
WorkingDirectory=/usr/local/sleep-manager
Environment=FLASK_APP=sleep-manager
Environment=FLASK_ENV=production
ExecStart=/usr/local/sleep-manager/venv/bin/python -m flask run --host=0.0.0.0 --port=51339
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/local/sleep-manager

[Install]
WantedBy=multi-user.target 