# systemd/sleep-manager.service
[Unit]
Description=Sleep Manager Service
After=network.target network-online.target multi-user.target
Wants=network.target network-online.target

[Service]
Type=simple
User=sleep-manager
Group=sleep-manager
WorkingDirectory=/usr/lib/sleep-manager
RuntimeDirectory=sleep-manager
RuntimeDirectoryMode=0750
ExecStart=/usr/lib/sleep-manager/venv/bin/gunicorn \
    --bind unix:/run/sleep-manager/sleep-manager.sock \
    --workers 2 \
    --timeout 30 \
    --access-logfile - \
    "sleep_manager:create_app()"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=control-group
TimeoutStopSec=10

# Security settings
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/lib/sleep-manager /run/sleep-manager

[Install]
WantedBy=multi-user.target
