# debian/postinst
#!/bin/sh
set -e

if ! id sleep-manager >/dev/null 2>&1; then
    useradd --system --user-group --shell /bin/false --home-dir /usr/lib/sleep-manager sleep-manager
fi

version_arg="$(dpkg-query -W -f='${Version}' sleep-manager 2>/dev/null || true)"
if [ -z "$version_arg" ]; then
    version_arg="${2:-}"
fi
version_arg="${version_arg#*:}"
version_arg="${version_arg%%-*}"
if [ -n "$version_arg" ]; then
    export HATCH_VCS_VERSION="$version_arg"
    export SETUPTOOLS_SCM_PRETEND_VERSION="$version_arg"
fi

mkdir -p /etc/sleep-manager
chown -R sleep-manager:sleep-manager /usr/lib/sleep-manager
chmod 755 /usr/lib/sleep-manager

needs_reenable=0
for unit in sleep-manager.service sleep-manager-sleeper.service sleep-manager-waker.service sleep-manager-delay.service; do
    if [ -d "/lib/systemd/system/$unit" ] && [ -f "/lib/systemd/system/$unit/$unit" ]; then
        mv "/lib/systemd/system/$unit/$unit" "/lib/systemd/system/$unit.tmp"
        rm -rf "/lib/systemd/system/$unit"
        mv "/lib/systemd/system/$unit.tmp" "/lib/systemd/system/$unit"
        needs_reenable=1
    fi
    if [ -f "/etc/systemd/system/$unit" ]; then
        if grep -q "/usr/local/sleep-manager" "/etc/systemd/system/$unit"; then
            mv "/etc/systemd/system/$unit" "/etc/systemd/system/${unit}.bak"
            needs_reenable=1
        fi
    fi
    if [ -d "/etc/systemd/system/${unit}.d" ]; then
        if grep -R -q "/usr/local/sleep-manager" "/etc/systemd/system/${unit}.d" 2>/dev/null; then
            mv "/etc/systemd/system/${unit}.d" "/etc/systemd/system/${unit}.d.bak"
            needs_reenable=1
        fi
    fi
done

for unit in sleep-manager-sleeper.service sleep-manager-waker.service; do
    if systemctl is-active "$unit" >/dev/null 2>&1; then
        systemctl stop "$unit" >/dev/null 2>&1 || true
        needs_reenable=1
    fi
    if systemctl is-enabled "$unit" >/dev/null 2>&1; then
        systemctl disable "$unit" >/dev/null 2>&1 || true
        needs_reenable=1
    fi
    if [ -f "/etc/systemd/system/$unit" ]; then
        rm -f "/etc/systemd/system/$unit"
        needs_reenable=1
    fi
done

example_config="/usr/lib/sleep-manager/config/sleep-manager-config.toml.example"
if [ -f /etc/sleep-manager/sleep-manager-config.json ]; then
    if [ ! -f /etc/sleep-manager/sleep-manager-config.toml ]; then
        tmp_config="$(mktemp /tmp/sleep-manager-config.toml.XXXXXX)"
        if /usr/lib/sleep-manager/scripts/migrate_config.py \
            /etc/sleep-manager/sleep-manager-config.json \
            "$tmp_config"; then
            if [ -s "$tmp_config" ]; then
                if command -v ucf >/dev/null 2>&1; then
                    ucf --debconf-ok "$tmp_config" /etc/sleep-manager/sleep-manager-config.toml
                else
                    mv "$tmp_config" /etc/sleep-manager/sleep-manager-config.toml
                fi
            else
                echo "WARNING: Config migration produced an empty TOML file; leaving config unset."
            fi
        else
            echo "WARNING: Config migration failed; leaving config unset."
        fi
        rm -f "$tmp_config"
    else
        is_example=0
        if [ -f "$example_config" ] && \
           cmp -s /etc/sleep-manager/sleep-manager-config.toml "$example_config"; then
            is_example=1
        elif grep -q "your-secure-api-key-here" /etc/sleep-manager/sleep-manager-config.toml 2>/dev/null; then
            is_example=1
        fi
    fi
    if [ "${is_example:-0}" -eq 1 ]; then
        rm -f /etc/sleep-manager/sleep-manager-config.toml
        tmp_config="$(mktemp /tmp/sleep-manager-config.toml.XXXXXX)"
        if /usr/lib/sleep-manager/scripts/migrate_config.py \
            /etc/sleep-manager/sleep-manager-config.json \
            "$tmp_config"; then
            if [ -s "$tmp_config" ]; then
                if command -v ucf >/dev/null 2>&1; then
                    ucf --debconf-ok "$tmp_config" /etc/sleep-manager/sleep-manager-config.toml
                else
                    mv "$tmp_config" /etc/sleep-manager/sleep-manager-config.toml
                fi
            else
                echo "WARNING: Config migration produced an empty TOML file; leaving config untouched."
            fi
        else
            echo "WARNING: Config migration failed; leaving config untouched."
        fi
        rm -f "$tmp_config"
    fi
elif [ -f /etc/sleep-manager/sleep-manager-config.toml ]; then
    if command -v ucf >/dev/null 2>&1; then
        tmp_config="$(mktemp /tmp/sleep-manager-config.toml.XXXXXX)"
        if /usr/lib/sleep-manager/scripts/migrate_toml_config.py \
            /etc/sleep-manager/sleep-manager-config.toml \
            --dest "$tmp_config"; then
            if [ -s "$tmp_config" ]; then
                ucf --debconf-ok "$tmp_config" /etc/sleep-manager/sleep-manager-config.toml
            else
                echo "WARNING: TOML config migration produced an empty file; leaving config untouched."
            fi
        else
            status=$?
            if [ "$status" -ne 2 ]; then
                echo "WARNING: TOML config migration failed; leaving config untouched."
            fi
        fi
        rm -f "$tmp_config"
    else
        echo "WARNING: ucf not available; skipping TOML config migration."
    fi
elif [ ! -f /etc/sleep-manager/sleep-manager-config.toml ] && [ -f "$example_config" ]; then
    cp "$example_config" /etc/sleep-manager/sleep-manager-config.toml
fi
if [ -f /etc/sleep-manager/sleep-manager-config.toml ]; then
    chown root:sleep-manager /etc/sleep-manager/sleep-manager-config.toml
    chmod 0640 /etc/sleep-manager/sleep-manager-config.toml
fi

if [ -f /etc/sudoers.d/sleep-manager ]; then
    chown root:root /etc/sudoers.d/sleep-manager
    chmod 0440 /etc/sudoers.d/sleep-manager
fi

if [ ! -d /usr/lib/sleep-manager/venv ]; then
    python3 -m venv /usr/lib/sleep-manager/venv
fi

chown -R sleep-manager:sleep-manager /usr/lib/sleep-manager/venv

/usr/lib/sleep-manager/venv/bin/pip install --upgrade pip
/usr/lib/sleep-manager/venv/bin/pip install --no-deps -r /usr/lib/sleep-manager/requirements.txt
/usr/lib/sleep-manager/venv/bin/pip install --upgrade hatchling hatch-vcs
/usr/lib/sleep-manager/venv/bin/pip install --no-deps --no-build-isolation /usr/lib/sleep-manager

# Install Homebridge plugin if npm is available (Homebridge is optional)
plugin_src="/usr/lib/sleep-manager/homebridge-sleep-manager"
hb_storage="/var/lib/homebridge"
if [ -d "$plugin_src" ] && [ -x /opt/homebridge/bin/homebridge ] && [ -x /opt/homebridge/bin/npm ] && [ -x /opt/homebridge/bin/node ]; then
    export PATH="/opt/homebridge/bin:$PATH"
    echo "Installing homebridge-sleep-manager plugin..."
    /opt/homebridge/bin/npm install --prefix "$hb_storage" "$plugin_src" \
        || echo "WARNING: Failed to install homebridge-sleep-manager; install manually with: env PATH=/opt/homebridge/bin:\$PATH npm install --prefix $hb_storage $plugin_src"
fi

if [ "$needs_reenable" -eq 1 ]; then
    systemctl daemon-reload >/dev/null 2>&1 || true
    for unit in sleep-manager.service sleep-manager-delay.service; do
        if systemctl is-enabled "$unit" >/dev/null 2>&1; then
            systemctl reenable "$unit" >/dev/null 2>&1 || true
        fi
    done
fi

# Clean up unit files from old rsync-installed location
rm -f /usr/lib/sleep-manager/systemd/sleep-manager.service
rm -f /usr/lib/sleep-manager/systemd/sleep-manager-delay.service

systemctl daemon-reload >/dev/null 2>&1 || true

# Enable nginx site
SITES_ENABLED=/etc/nginx/sites-enabled
SITES_AVAILABLE=/etc/nginx/sites-available
if [ ! -L "$SITES_ENABLED/sleep-manager" ]; then
    ln -s "$SITES_AVAILABLE/sleep-manager" "$SITES_ENABLED/sleep-manager"
fi
if nginx -t -q 2>/dev/null && systemctl is-active --quiet nginx; then
    systemctl reload nginx
fi

#DEBHELPER#
