# debian/postinst
#!/bin/sh
set -e

if ! id sleep-manager >/dev/null 2>&1; then
    useradd --system --user-group --shell /bin/false --home-dir /usr/lib/sleep-manager sleep-manager
fi

version_arg="$(dpkg-query -W -f='${Version}' sleep-manager 2>/dev/null || true)"
if [ -z "$version_arg" ]; then
    version_arg="${2:-}"
fi
version_arg="${version_arg#*:}"
version_arg="${version_arg%%-*}"
if [ -n "$version_arg" ]; then
    export HATCH_VCS_VERSION="$version_arg"
    export SETUPTOOLS_SCM_PRETEND_VERSION="$version_arg"
fi

mkdir -p /etc/sleep-manager
chown -R sleep-manager:sleep-manager /usr/lib/sleep-manager
chmod 755 /usr/lib/sleep-manager

needs_reenable=0
for unit in sleep-manager-sleeper.service sleep-manager-waker.service sleep-manager-delay.service; do
    if [ -d "/lib/systemd/system/$unit" ] && [ -f "/lib/systemd/system/$unit/$unit" ]; then
        mv "/lib/systemd/system/$unit/$unit" "/lib/systemd/system/$unit.tmp"
        rm -rf "/lib/systemd/system/$unit"
        mv "/lib/systemd/system/$unit.tmp" "/lib/systemd/system/$unit"
        needs_reenable=1
    fi
    if [ -f "/etc/systemd/system/$unit" ]; then
        if grep -q "/usr/local/sleep-manager" "/etc/systemd/system/$unit"; then
            mv "/etc/systemd/system/$unit" "/etc/systemd/system/${unit}.bak"
            needs_reenable=1
        fi
    fi
    if [ -d "/etc/systemd/system/${unit}.d" ]; then
        if grep -R -q "/usr/local/sleep-manager" "/etc/systemd/system/${unit}.d" 2>/dev/null; then
            mv "/etc/systemd/system/${unit}.d" "/etc/systemd/system/${unit}.d.bak"
            needs_reenable=1
        fi
    fi
done

if [ -f /etc/sleep-manager/sleep-manager-config.json ]; then
    chown root:sleep-manager /etc/sleep-manager/sleep-manager-config.json
    chmod 0640 /etc/sleep-manager/sleep-manager-config.json
fi

if [ ! -d /usr/lib/sleep-manager/venv ]; then
    python3 -m venv /usr/lib/sleep-manager/venv
fi

chown -R sleep-manager:sleep-manager /usr/lib/sleep-manager/venv

/usr/lib/sleep-manager/venv/bin/pip install --upgrade pip setuptools wheel
/usr/lib/sleep-manager/venv/bin/pip install --upgrade hatchling hatch-vcs
/usr/lib/sleep-manager/venv/bin/pip install --upgrade --no-build-isolation /usr/lib/sleep-manager

if [ "$needs_reenable" -eq 1 ]; then
    systemctl daemon-reload >/dev/null 2>&1 || true
    for unit in sleep-manager-sleeper.service sleep-manager-waker.service sleep-manager-delay.service; do
        if systemctl is-enabled "$unit" >/dev/null 2>&1; then
            systemctl reenable "$unit" >/dev/null 2>&1 || true
        fi
    done
fi
systemctl daemon-reload >/dev/null 2>&1 || true

exit 0
