# debian/sleep-manager.sleeper.service
[Unit]
Description=Sleep Manager Sleeper Service
After=network.target network-online.target multi-user.target
Wants=network.target network-online.target

[Service]
Type=simple
User=sleep-manager
Group=sleep-manager
WorkingDirectory=/usr/lib/sleep-manager
# Use the Flask app factory pattern
Environment=FLASK_APP=sleep_manager:create_app
Environment=FLASK_ENV=production
ExecStart=/usr/lib/sleep-manager/venv/bin/python -m flask run --host=0.0.0.0 --port=51339
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
KillMode=control-group
TimeoutStopSec=10

# Security settings
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/usr/lib/sleep-manager

[Install]
WantedBy=multi-user.target
